# packages/atlas-hub-<version>/supervisord.conf
[unix_http_server]
file=/usr/lib/atlas-hub/supervisor.sock   ; the path to the socket file

[supervisord]
logfile=/var/log/atlas-hub/supervisord.log
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=/usr/lib/atlas-hub/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=false               ; start in foreground if true; default false
silent=false                 ; no logs to stdout if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200
;umask=022                   ; process file creation umask; default 022
;user=supervisord            ; setuid to this UNIX account at startup; recommended if root
;identifier=supervisor       ; supervisord identifier, default is 'supervisor'
;directory=/tmp              ; default is not to cd during start
;nocleanup=true              ; don't clean up tempfiles at start; default false
;childlogdir=/tmp            ; 'AUTO' child log dir, default $TEMP
;environment=KEY="value"     ; key value pairs to add to environment
;strip_ansi=false            ; strip ansi escape codes in logs; def. false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///usr/lib/atlas-hub/supervisor.sock ; use a unix:// URL  for a unix socket

[program:web]
directory=/usr/lib/atlas-hub/app
command=/usr/lib/atlas-hub/app/.venv/bin/gunicorn --worker-class=gevent --workers 1 --threads 30 --timeout 999999999 --access-logfile /var/log/atlas-hub/access.log --error-logfile /var/log/atlas-hub/error.log --capture-output --bind  unix:web.sock --umask 007 web:app
autostart=true
autorestart=true
user=www-data

[program:runner]
directory=/usr/lib/atlas-hub/app
command=/usr/lib/atlas-hub/app/.venv/bin/gunicorn --worker-class=gevent --workers 1 --threads 30 --timeout 999999999 --access-logfile /var/log/atlas-hub/access.log --error-logfile /var/log/atlas-hub/error.log --capture-output --bind  unix:runner.sock --umask 007 runner:app
autostart=true
autorestart=true
user=www-data

[program:scheduler]
directory=/usr/lib/atlas-hub/app
command=/usr/lib/atlas-hub/app/.venv/bin/gunicorn --worker-class=gevent --workers 1 --threads 30 --timeout 999999999 --access-logfile /var/log/atlas-hub/access.log --error-logfile /var/log/atlas-hub/error.log --capture-output --bind  unix:scheduler.sock --umask 007 scheduler:app
autostart=true
autorestart=true
user=www-data
